# roles/pihole/files/pihole-up.timer
[Timer]
OnCalendar=Sat 05:00
Persistent=true

[Install]
WantedBy=pihole-FTL.service
