# roles/monitor/files/src/units/php-fpm.service.d/override.conf
[Service]
ReadWritePaths = /etc/icingaweb2
