# roles/cloud/files/container/src/units/cloud-cron.timer
[Unit]
Description=Cloud cron
BindsTo=php-fpm.service

[Timer]
OnUnitActiveSec=5min
OnBootSec=5min

[Install]
WantedBy=php-fpm.service
