# roles/cloud/files/src/units/cloud-previews.timer
[Unit]
BindsTo=php-fpm.service

[Timer]
OnUnitActiveSec=10min

[Install]
WantedBy=php-fpm.service
