# roles/cloud/files/src/units/cloud-previews.timer
[Unit]

[Timer]
OnCalendar=Sun

[Install]
WantedBy=php-fpm.service
