# roles/cloud/files/src/units/cloud-previews.timer
[Unit]
Description=Cloud previews
BindsTo=php-fpm.service

[Timer]
OnUnitActiveSec=10min
OnBootSec=10min

[Install]
WantedBy=php-fpm.service
