# roles/cloud/files/src/units/cloud-cron.timer
[Unit]
BindsTo=php-fpm.service

[Timer]
OnUnitActiveSec=5min

[Install]
WantedBy=php-fpm.service
