# roles/monitor/files/src/units/container.target
[Unit]
Requires=php-fpm.service caddy.service icinga2.service icingadb.service
After=php-fpm.service caddy.service icinga2.service icingadb.service
