# roles/common/files/saturday.timer
[Timer]
OnCalendar=
OnCalendar=Sat 05:00
