# roles/container/files/units/image-prune.timer
[Timer]
OnCalendar = 01:00
Persistent = true

[Install]
WantedBy = podman-auto-update.timer
