# roles/container/files/units/image-build@.timer
[Unit]
Description=Image build

[Timer]
OnCalendar = Fri 23:00
Persistent = true

[Install]
WantedBy = podman-auto-update.timer
