# roles/container/files/units/prod-builders.slice
[Unit]
Description=Image build queue

[Slice]
# Excess jobs will wait in the queue until a slot opens.
ConcurrencySoftMax=2
