# roles/dyndns/files/strato-dyndns.timer
[Unit]
Description=Periodic Strato DynDNS update

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
