# roles/dyndns/files/strato-dyndns.service
[Unit]
Description=Update Strato DynDNS AAAA record
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=STRATO_DYNDNS_USER=dm-poepperl.de
Environment=STRATO_DYNDNS_HOST=dm-poepperl.de
LoadCredentialEncrypted=strato-password:/etc/strato-dyndns.cred
ExecStart=/usr/bin/python3 /usr/local/bin/strato-dyndns.py
