# roles/backup_base/files/syncoid@.timer
[Timer]
OnCalendar = *:10
Persistent = true
RandomizedDelaySec = 2m

[Install]
WantedBy = timers.target
